# agent/bsm-agent.service
[Unit]
Description=BSM Community privileged allowlisted agent
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/bsm-community/current/agent/bsm_agent.py
Restart=on-failure
RestartSec=3
UMask=0077
PrivateTmp=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
